# archive/baseball-ai-media-backup/systemd_config/bbdata-deep@.service
[Unit]
Description=BBData Deep (CSR/JS) Collector for %i
After=network-online.target

[Service]
Type=simple
User=mizu
Group=mizu
WorkingDirectory=/home/mizu/baseball-ai-media
Environment=TZ=Asia/Tokyo
Environment=VIRTUAL_ENV=/home/mizu/baseball-ai-media/.venv
Environment=PATH=/home/mizu/baseball-ai-media/.venv/bin:/usr/local/bin:/usr/bin:/bin
# 安全弁
Environment=BB_MIN_INTERVAL_SEC=15
Environment=BB_MAX_PAGES_PER_RUN=300
Environment=BB_MAX_CONCURRENT=1
Environment=BB_BACKOFF_429_SEC=300
ExecStart=/home/mizu/baseball-ai-media/.venv/bin/python scripts/bbdata_deep_collect.py --date %i --targets vducp,conditions --changed-only --discord
Nice=10
Restart=no
StandardOutput=journal
StandardError=journal
TimeoutStartSec=7200

[Install]
WantedBy=multi-user.target