# archive/deploy/baseball-live-api.service
[Unit]
Description=Baseball AI Media - Live Prediction API Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=baseball
Group=baseball
WorkingDirectory=/opt/baseball-ai-media
ExecStart=/usr/bin/npx tsx scripts/live-server.ts
Environment=NODE_ENV=production
Environment=LIVE_PORT=8787
Environment=DB_PATH=./data/baseball_live.db
Environment=MODELS_PATH=./models
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=baseball-live-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/baseball-ai-media/data /opt/baseball-ai-media/logs /opt/baseball-ai-media/models

# Resource limits for SSE connections
LimitNOFILE=65536
MemoryMax=4G
TasksMax=2048

[Install]
WantedBy=multi-user.target