# archive/systemd_config/bbdata-roster.service
[Unit]
Description=BaseballData.jp Player Roster Daily Collection
After=network.target

[Service]
Type=oneshot
User=mizu
Group=mizu
WorkingDirectory=/home/mizu/baseball-ai-media
Environment=VIRTUAL_ENV=/home/mizu/baseball-ai-media/.venv
Environment=PATH=/home/mizu/baseball-ai-media/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/mizu/baseball-ai-media/.venv/bin/python scripts/bbdata_player_roster.py --date $(date +%%Y-%%m-%%d) --discord
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target