# archive/systemd_config/bbdata-gamelog.service
[Unit]
Description=BaseballData.jp Player Gamelog Daily Collection
After=network.target bbdata-stats.service

[Service]
Type=oneshot
User=mizu
Group=mizu
WorkingDirectory=/home/mizu/baseball-ai-media
Environment=VIRTUAL_ENV=/home/mizu/baseball-ai-media/.venv
Environment=PATH=/home/mizu/baseball-ai-media/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/mizu/baseball-ai-media/.venv/bin/python scripts/bbdata_player_gamelog.py --date $(date +%%Y-%%m-%%d) --discord
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target