# archive/systemd_config/npb-step2.timer
[Unit]
Description=Run NPB Step2 every 60 minutes (polite)

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
RandomizedDelaySec=120
Unit=npb-step2.service

[Install]
WantedBy=timers.target