# archive/systemd_config/npb-step3-live.service
[Unit]
Description=NPB Step3 Pitch Log Collector (polite, continuous)
After=network-online.target

[Service]
EnvironmentFile=/etc/default/npb-scraper
WorkingDirectory=%E{BASE_DIR}
Environment=TZ=%E{TZ}
ExecStart=/bin/bash -lc 'source %E{VENV}/bin/activate && python scripts/step_3_pitchlog_fetcher_polite.py'
Restart=always
RestartSec=30s
Nice=10
IOSchedulingClass=best-effort
# 安全に落とす（Ctrl+C相当の SIGINT → 60s 後に SIGKILL）
KillSignal=SIGINT
TimeoutStopSec=60