# archive/systemd_config/npb-step1.timer
[Unit]
Description=Run NPB Step1 every 15 minutes (polite)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=60
Unit=npb-step1.service

[Install]
WantedBy=timers.target